# scripts/com.anna.spending-calendar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Launchd agent — runs refresh_local.sh daily at 9am local time.

  One-time setup:
    cp scripts/com.anna.spending-calendar.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.anna.spending-calendar.plist

  REQUIRED if the repo lives under ~/Desktop (or ~/Documents): macOS TCC blocks
  launchd from executing scripts in those folders, so the job dies with exit 126
  ("Operation not permitted"). Grant Full Disk Access to /bin/bash:
    System Settings → Privacy & Security → Full Disk Access → + → /bin/bash
  then unload + load the agent again. (The `bash -l` flag does NOT fix this.)

  To test it immediately (without waiting for 9am):
    launchctl start com.anna.spending-calendar

  To stop / unload:
    launchctl unload ~/Library/LaunchAgents/com.anna.spending-calendar.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.anna.spending-calendar</string>

  <key>Program</key>
  <string>/bin/bash</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>/Users/anna/Desktop/Home/AI-Projects/calendar/scripts/refresh_local.sh</string>
  </array>

  <!-- Run daily at 9:00am local time (America/Denver) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- If the Mac was asleep at 9am, run as soon as it wakes -->
  <key>StartInterval</key>
  <integer>86400</integer>

  <key>StandardOutPath</key>
  <string>/Users/anna/Desktop/Home/AI-Projects/calendar/logs/refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/anna/Desktop/Home/AI-Projects/calendar/logs/refresh.log</string>

  <!-- Don't run immediately on load — wait for the scheduled time -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
